# systemd/tresk-deep-scan.timer
[Unit]
Description=Tresk - Daily Deep Scan Timer
Documentation=https://github.com/yourusername/tresk

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
