# systemd/tresk.service
[Unit]
Description=Tresk - VPS Security Monitor
Documentation=https://github.com/yourusername/tresk
After=network.target auditd.service
Wants=auditd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/tresk

# Main monitoring script
ExecStart=/opt/tresk/bin/monitor.sh monitor
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
CPUQuota=10%
MemoryLimit=100M
TasksMax=50

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/var/log/tresk /opt/tresk

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tresk

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="CONFIG_FILE=/etc/tresk/config.conf"

[Install]
WantedBy=multi-user.target
