# systemd/tresk-summary.timer
[Unit]
Description=Tresk - Daily Summary Timer
Documentation=https://github.com/yourusername/tresk

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
